# bashLessons/site-backup.service
[Unit]
Description=Backup BindAndDNS

[Service]
Type=simple
ExecStart=/home/gweep/Documets/site-backup.sh

[Install]
WantedBy=multi-user.target
